# backup_content/crontab
MAILTO="tuatini@theenk.io"
@midnight root /root/backup_wordpress_data.sh
